# dist/gentoo/app-pda/synce-dtpt/files/init.dtpt
#!/sbin/runscript

PIDFILE="/var/run/dtpt.pid"

start() {
   	ebegin "Starting DTPT daemon"
	start-stop-daemon --start --quiet --background --make-pidfile --pidfile "$PIDFILE" --exec `which python` -- /usr/share/dtpt/DTPT.py
	eend $?
	eend 0
}

stop() {
	ebegin "Shutting down DTPT daemon"
	start-stop-daemon --stop --quiet --pidfile "$PIDFILE"
	eend $?
	rm -f "$PIDFILE"
	eend 0
}
