# gvfs/daemon/synce.mount.in
[Mount]
Type=synce
Exec=@libexecdir@/gvfsd-synce
AutoMount=false
Scheme=synce
